# deploy/homecontrol.service
[Unit]
Description="Homecontrol app"
After=network.target

[Service]
User=ahs
RuntimeDirectory=homecontrol
WorkingDirectory=/home/ahs/src/homecontrol
ExecStart=/home/ahs/.virtualenvs/homecontrol/venv/bin/python run.py

[Install]
WantedBy=multi-user.target
